# packages/repo/w/wpa_supplicant/.install/install.sh
systemctl enable wpa_supplicant
update-desktop-database -q
