# packages/repo/s/systemd/.install/install.sh
systemd-machine-id-setup

systemctl preset-all
